# vertretung.service
[Unit]
Description=Den Vertretungsplan als Benachrichtung anzeigen

[Service]
ExecStart=/home/&0/bin/check_vertretungsplan
